# deploy/pulpo-ocr-worker.service
[Unit]
Description=Pulpo OCR Worker
After=network.target
Wants=network.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
Environment=OCR_RUN_MODE=loop
Environment=OCR_LOOP_INTERVAL=30
Environment=OCR_MAX_CONCURRENCY=2
Environment=OCR_MAX_RETRIES=3
Environment=TESSERACT_CMD=tesseract
Environment=TESSERACT_LANG=spa
WorkingDirectory=/home/nictes/workspace/nictes1/pulpo
ExecStart=/home/nictes/workspace/nictes1/pulpo/.venv/bin/python -m services.ocr_worker
Restart=always
RestartSec=10
User=nictes
Group=nictes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pulpo-ocr-worker

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/nictes/workspace/nictes1/pulpo

[Install]
WantedBy=multi-user.target
